# pub/Linux/adaptorsos/lib/systemd/system/adaptorsos.service
[Unit]
Description=Adaptor SOS

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adaptorsos/adaptorsos.sh
