# pub/Linux/adaptorsos/DEBIAN/prerm
#!/bin/bash

if ( systemctl is-active --quiet multiplug.service )
then
systemctl disable multiplug.service
fi